# vampigotchi.service
[Unit]
Description=VampiGotchi BLE Monitoring Service
After=network.target bluetooth.target

[Service]
Type=simple
User=root
WorkingDirectory=/path/to/VAMPIGOTCHI
ExecStart=/usr/bin/python3 /path/to/Scrpit.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
